# deploy/systemd/ch-live-catchup.timer
[Unit]
Description=Rates Engine — CH live fan-out catch-up timer (every 10 min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
