# deploy/systemd/ch-supply-refresh.service
[Unit]
Description=Rates Engine — ClickHouse per-token supply refresh (stellar.token_supply, ADR-0034)
After=clickhouse-server.service network-online.target
Wants=clickhouse-server.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ch-supply-refresh.sh
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
