# deploy/systemd/sla-probe.service
# systemd unit: SLA-probe periodic runner.
#
# Drives `ratesengine-sla-probe` against the configured API base URL,
# producing per-endpoint p50 / p95 / p99 latency, freshness, and
# availability values, with a pass/fail verdict against the
# RFP-stated SLA targets. Companion to sla-probe.timer which
# schedules the cadence.
#
# Operator wiring:
#
#   sudo cp deploy/systemd/sla-probe.{service,timer} /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sla-probe.timer
#
# Configuration (override via /etc/default/sla-probe):
#
#   BASE_URL                    API base URL.       (default: https://api.ratesengine.net/v1)
#   DURATION                    Per-run load duration.   (default: 30s)
#   CONCURRENCY                 Concurrent workers.      (default: 4)
#   PAIRS                       One -pair flag per pair. (default: -pair native,fiat:USD)
#   REPORT_FORMAT               text or json.            (default: json — easy to scrape from journal)
#   RATESENGINE_PROBE_API_KEY   Vault-minted load-test API key. REQUIRED — without
#                               it the probe hits the anonymous-tier rate limit
#                               (60 req/min) and reads availability < 0.1% on every
#                               non-/healthz endpoint. Mint same class as the k6
#                               weekly's RATESENGINE_LOAD_API_KEY. Picked up by the
#                               binary from this env var directly so the key never
#                               appears on the ExecStart command line.
#   TEXTFILE_OUTPUT             Path under node_exporter's textfile_collector dir.
#                               Empty = no metrics emit (the default — operators
#                               opt in once they've installed node_exporter).
#                               Example:
#                                 /var/lib/node_exporter/textfile_collector/sla_probe.prom
#   EXTRA_FLAGS                 Extra -flag overrides.   (default: empty)
#
# Exit code 0 = pass; 1 = at least one SLA violated. The non-zero
# exit on SLA breach causes systemd to log the unit as failed,
# which surfaces via the existing
# `ratesengine_sla_probe_unit_failed` Prometheus alert (planned
# follow-up that pairs with the textfile-collector integration).

[Unit]
Description=Rates Engine — SLA probe (periodic per-endpoint p95/p99/freshness check)
Documentation=https://github.com/RatesEngine/rates-engine/blob/main/docs/operations/sla-probe.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/sla-probe

# Sensible defaults. Operators override via /etc/default/sla-probe.
Environment=BASE_URL=https://api.ratesengine.net/v1
Environment=DURATION=30s
Environment=CONCURRENCY=4
# F-1205 (codex audit-2026-05-12): the value MUST be quoted —
# systemd's Environment= parser splits on whitespace UNLESS the
# value is in quotes, so the unquoted form was being parsed as
# two assignments (`PAIRS=-pair` plus a stray `native,fiat:USD`),
# leaving the probe to fail with `flag needs an argument: -pair`.
Environment="PAIRS=-pair native,fiat:USD"
Environment=REPORT_FORMAT=json
Environment=TEXTFILE_OUTPUT=
Environment=EXTRA_FLAGS=

ExecStart=/bin/sh -c '/usr/local/bin/ratesengine-sla-probe \
  -base-url ${BASE_URL} \
  -duration ${DURATION} \
  -concurrency ${CONCURRENCY} \
  -report-format ${REPORT_FORMAT} \
  $(if [ -n "${TEXTFILE_OUTPUT}" ]; then echo "-textfile-output ${TEXTFILE_OUTPUT}"; fi) \
  ${PAIRS} \
  ${EXTRA_FLAGS}'

# Network-only — no filesystem writes. Runs as the application user.
User=ratesengine
Group=ratesengine

# Hardening — read-only filesystem except for the textfile-collector
# dir when TEXTFILE_OUTPUT is set. Operators that change the path
# below (or set it to a non-default location) must update
# ReadWritePaths to match — otherwise the atomic-write rename fails
# under ProtectSystem=strict.
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector

# A run is bounded by DURATION + a small startup tax. Cap at 5x the
# default duration so a stuck probe can't pin the unit beyond one
# timer window.
TimeoutStartSec=2min30s

[Install]
WantedBy=multi-user.target
