# configs/ansible/roles/archival-node/templates/systemd/galexie.service.j2
[Unit]
Description=Stellar Galexie — ledger-meta exporter
Documentation=https://github.com/stellar/stellar-galexie
After=network-online.target{% if run_minio %} minio.service{% endif %}

Wants=network-online.target

[Service]
Type=simple
User=galexie
Group=galexie
WorkingDirectory=/var/lib/galexie
EnvironmentFile=-/etc/default/galexie
# Wrapper handles:
#  * flag is --config-file (not --config)
#  * `append` needs an explicit --start ledger; we discover the
#    current network tip from stellar-core on each service start
#    so galexie begins exporting from "now" on a fresh deployment,
#    or resumes from its last-exported-ledger+1 if MinIO already
#    has a cursor for us (the wrapper probes for it).
ExecStart=/usr/local/bin/galexie-append.sh
ExecStop=/bin/kill -TERM $MAINPID

LimitNOFILE=1048576
LimitNPROC=1048576

Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
