# configs/ansible/roles/archival-node/templates/systemd/ratesengine-aggregator.service.j2
[Unit]
Description=Rates Engine aggregator (VWAP/TWAP CAGG refresh + Redis cache)
Documentation=https://github.com/RatesEngine/rates-engine/blob/main/docs/architecture/aggregation-plan.md
After=network-online.target postgresql@{{ postgres_version }}-main.service redis-server.service
Wants=network-online.target
Requires=postgresql@{{ postgres_version }}-main.service

[Service]
Type=simple
# Run as root for now — same convention as galexie + indexer.
# Tighten to a dedicated `ratesengine` system user once 04-users.yml
# grows one.
User=root
Group=root
EnvironmentFile=/etc/default/ratesengine-ops
ExecStart=/usr/local/bin/ratesengine-aggregator -config /etc/ratesengine.toml
ExecStop=/bin/kill -TERM $MAINPID

LimitNOFILE=1048576
Restart=on-failure
RestartSec=10s

# Sandboxing — same defaults as the indexer.
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
