# configs/ansible/roles/archival-node/templates/systemd/ratesengine-api.service.j2
[Unit]
Description=Rates Engine — REST + SSE API server
After=network-online.target postgresql@{{ postgres_version }}-main.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/ratesengine
EnvironmentFile=/etc/default/ratesengine
ExecStart=/usr/local/bin/ratesengine-api -config /etc/ratesengine.toml
ExecStop=/bin/kill -TERM $MAINPID

LimitNOFILE=1048576
Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
