# configs/ansible/roles/archival-node/templates/systemd/supply-snapshot.service.j2
# systemd unit: supply-snapshot writer.
#
# Rendered by Ansible. Reflects r1's EFFECTIVE running config — the
# base unit plus the override.conf drop-in (User=root + the
# /etc/default/ratesengine EnvironmentFile) is folded directly into
# this template.
#
# Per ADR-0011 + the supply-snapshot.timer cadence. Runs once on
# each timer tick (oneshot semantics); the timer file at
# supply-snapshot.timer schedules the cadence.
#
# Configuration (override via /etc/default/supply-snapshot):
#   CONFIG_PATH      Default /etc/ratesengine.toml
#   ASSET            Default native (only XLM at v1).
#   TEXTFILE_OUTPUT  Path under node_exporter's textfile_collector
#                    dir. Empty = no metrics emit.
#   EXTRA_FLAGS      Optional flags appended to the command line.

[Unit]
Description=Rates Engine — supply-snapshot writer (daily cron per ADR-0011)
Documentation=https://github.com/RatesEngine/rates-engine/blob/main/docs/operations/supply-snapshot.md
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/supply-snapshot
EnvironmentFile=-/etc/default/ratesengine

# Sensible defaults. Operators override via /etc/default/supply-snapshot.
Environment=CONFIG_PATH=/etc/ratesengine.toml
Environment=ASSET=native
Environment=TEXTFILE_OUTPUT=
Environment=EXTRA_FLAGS=

ExecStart=/bin/sh -c '/usr/local/bin/ratesengine-ops supply snapshot \
  -config ${CONFIG_PATH} \
  -asset ${ASSET} \
  $(if [ -n "${TEXTFILE_OUTPUT}" ]; then echo "-textfile-output ${TEXTFILE_OUTPUT}"; fi) \
  ${EXTRA_FLAGS}'

# Runs as root for consistency with the ratesengine-* daemons (the
# in-repo unit set User=ratesengine but r1 has no such user yet —
# F-1205, audit-2026-05-12).
User=root
Group=root

# Hardening — config is read-only; textfile_collector dir is the
# only writable path. Operators that change TEXTFILE_OUTPUT to a
# non-default location must update ReadWritePaths to match.
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector

# A snapshot run is fast: one ListCursors + one Compute + one
# InsertSupply. Worst case is a slow Postgres — cap at 5 minutes
# so a wedged DB can't pin the unit.
TimeoutStartSec=5min

[Install]
WantedBy=multi-user.target
